# backend/hideandseek-backend.service
# systemd unit file – alternative to PM2
# Install: sudo cp hideandseek-backend.service /etc/systemd/system/
#          sudo systemctl enable --now hideandseek-backend
[Unit]
Description=Hide and Seek Backend
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/hideandseek/backend
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5s
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=DB_PATH=/var/data/hideandseek/hideandseek.db
Environment=FRONTEND_ORIGIN=https://yourdomain.com
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
